# services/gbrain-embeddings/gbrain-embeddings.service
[Unit]
Description=GBrain local embedding stack
Documentation=file:/home/sailorjoe6/Code/vllm-sleeper-proxy/services/gbrain-embeddings/README.md
Requires=docker.service
After=network-online.target docker.service nvidia-persistenced.service
Wants=network-online.target nvidia-persistenced.service

[Service]
Type=simple
WorkingDirectory=/home/sailorjoe6/Code/vllm-sleeper-proxy/services/gbrain-embeddings
Environment=HOME=/home/sailorjoe6
Environment=HF_CACHE_DIR=/home/sailorjoe6/.cache/huggingface
ExecStart=/usr/bin/docker compose up --remove-orphans
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
